# setup/systemctl/zkill-cron.service
[Unit]
Description=ZKill-Cron

[Service]
ExecStart=/usr/bin/node ./bin/cron.js
Restart=always
User=zkill
Group=zkill
Environment=PATH=/usr/bin:/usr/local/bin
EnvironmentFile=/home/zkill/Dropbox/zkilljs/.env
WorkingDirectory=/home/zkill/Dropbox/zkilljs

[Install]
WantedBy=multi-user.target
